# tweetloop.service
[Unit]
Description=TweetLoop App
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/danny/workspace/twitter-reviewer
ExecStart=/home/danny/.hermes/hermes-agent/venv/bin/python3 /home/danny/workspace/twitter-reviewer/app.py
Environment=PORT=7777
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
